# autostart/.config/autostart/flameshot.desktop
[Desktop Entry]
Type=Application
Name=Flameshot
NoDisplay=false
Exec=/usr/bin/flameshot &
Comment=Screenshot Application.
Categories=System;
